# pisignage.service
# PiSignage v0.8.0 - Service Systemd
# Service principal pour l'affichage digital
# Auteur: Claude Code
# Date: 22/09/2025

[Unit]
Description=PiSignage Digital Signage v0.8.0
Documentation=https://github.com/elkir0/Pi-Signage
After=network.target nginx.service php8.2-fpm.service
Wants=nginx.service php8.2-fpm.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=forking
User=pi
Group=pi
WorkingDirectory=/opt/pisignage

# Variables d'environnement
Environment=DISPLAY=:0
Environment=PISIGNAGE_DIR=/opt/pisignage
Environment=MEDIA_DIR=/opt/pisignage/media
Environment=LOG_DIR=/opt/pisignage/logs

# Commande de démarrage
ExecStartPre=/bin/bash -c 'if [ ! -f /tmp/.X0-lock ]; then startx & sleep 5; fi'
ExecStart=/opt/pisignage/scripts/pisignage-start.sh
ExecReload=/opt/pisignage/scripts/pisignage-reload.sh
ExecStop=/opt/pisignage/scripts/pisignage-stop.sh

# Gestion des erreurs et redémarrages
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=60
TimeoutStopSec=30

# Sécurité
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=false
PrivateTmp=false

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pisignage

[Install]
WantedBy=multi-user.target
Alias=pisignage.service