# deploy/systemd/pisignage-audio-apply.service
[Unit]
Description=Zaforge: applique la sortie audio persistée (HDMI/jack) au boot
# WirePlumber tourne en service USER (pas system) sur Trixie. On peut pas
# référencer wireplumber.service ici. À la place : attendre graphical.target
# (lightdm + session pi complète, qui lance wireplumber user). Le script apply
# retry 12x5s=1min si le sink n'est pas encore créé (TV branchée tard).
After=graphical.target

[Service]
Type=oneshot
User=pi
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
ExecStart=/opt/pisignage/scripts/audio-output-apply.sh
# Reste "active" après success (juste pour le status systemd clair).
RemainAfterExit=yes
# Retry au boot si la session pi n'était pas tout à fait prête.
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
