# templates/autostart/pisignage.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=PiSignage Desktop
Comment=Digital Signage Player for Raspberry Pi Desktop
Comment[fr]=Lecteur d'affichage numérique pour Raspberry Pi Desktop
GenericName=Digital Signage Player
GenericName[fr]=Lecteur d'affichage numérique
Keywords=signage;digital;display;multimedia;video;kiosk;
Keywords[fr]=affichage;numérique;écran;multimédia;vidéo;kiosque;

# Commande d'exécution
Exec=/opt/pisignage/pisignage-desktop/scripts/control/start.sh
TryExec=/opt/pisignage/pisignage-desktop/scripts/control/start.sh
Path=/opt/pisignage/pisignage-desktop

# Configuration d'affichage
Icon=/opt/pisignage/pisignage-desktop/web/assets/icon.png
StartupNotify=false
NoDisplay=false
Hidden=false

# Autostart configuration
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=10
X-KDE-autostart-after=panel
X-MATE-Autostart-Delay=10

# Categories
Categories=AudioVideo;Player;TV;
MimeType=video/mp4;video/avi;video/mkv;video/mov;video/webm;application/x-pisignage-playlist;

# Actions
Actions=Start;Stop;Restart;Status;

[Desktop Action Start]
Name=Start PiSignage
Name[fr]=Démarrer PiSignage
Exec=/opt/pisignage/pisignage-desktop/scripts/control/start.sh
Icon=/opt/pisignage/pisignage-desktop/web/assets/icon-start.png

[Desktop Action Stop]
Name=Stop PiSignage
Name[fr]=Arrêter PiSignage
Exec=/opt/pisignage/pisignage-desktop/scripts/control/stop.sh
Icon=/opt/pisignage/pisignage-desktop/web/assets/icon-stop.png

[Desktop Action Restart]
Name=Restart PiSignage
Name[fr]=Redémarrer PiSignage
Exec=/opt/pisignage/pisignage-desktop/scripts/control/restart.sh
Icon=/opt/pisignage/pisignage-desktop/web/assets/icon-restart.png

[Desktop Action Status]
Name=Show Status
Name[fr]=Afficher le statut
Exec=/opt/pisignage/pisignage-desktop/scripts/control/status.sh
Icon=/opt/pisignage/pisignage-desktop/web/assets/icon-status.png
Terminal=true