# systemd/pisignage.service
[Unit]
Description=PiSignage Digital Signage System
Documentation=file:///opt/pisignage/docs/README.md
After=network.target sound.target
Wants=network.target

[Service]
Type=forking
User=pi
Group=pi
WorkingDirectory=/opt/pisignage
Environment=DISPLAY=:0
Environment=HOME=/home/pi

# Scripts de démarrage et d'arrêt
ExecStartPre=/bin/bash -c 'mkdir -p /opt/pisignage/logs'
ExecStart=/opt/pisignage/player-control.sh start
ExecStop=/opt/pisignage/player-control.sh stop
ExecReload=/opt/pisignage/player-control.sh restart

# Redémarrage automatique en cas d'erreur
Restart=on-failure
RestartSec=10

# Timeout pour les opérations
TimeoutStartSec=60
TimeoutStopSec=30

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pisignage

# Sécurité
NoNewPrivileges=true
PrivateTmp=true

# Variables d'environnement supplémentaires
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=multi-user.target
Alias=pisignage