# templates/systemd/pisignage-web.service
[Unit]
Description=PiSignage Desktop Web Service
After=network.target nginx.service
Requires=nginx.service

[Service]
Type=simple
User=www-data
Group=www-data
ExecStart=/opt/pisignage/pisignage-desktop/scripts/control/web-service.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target