# old-version/config/pisignage-v4.service
[Unit]
Description=PiSignage v4.0 - Digital Signage Engine
Documentation=https://github.com/elkir0/Pi-Signage
After=graphical-session.target network.target sound.target
Wants=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=forking
User=pi
Group=video
Environment=DISPLAY=:0
Environment=HOME=/home/pi
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/opt/pisignage

# Script de démarrage principal
ExecStart=/opt/pisignage/scripts/vlc-v4-engine.sh start
ExecStop=/opt/pisignage/scripts/vlc-v4-engine.sh stop
ExecReload=/opt/pisignage/scripts/vlc-v4-engine.sh restart

# Gestion des processus
PIDFile=/opt/pisignage/run/vlc-engine.pid
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=30
TimeoutStopSec=15
RestartSec=5

# Politique de redémarrage robuste
Restart=always
RestartPreventExitStatus=0

# Permissions et sécurité
SupplementaryGroups=audio video render input
DeviceAllow=/dev/fb0 rw
DeviceAllow=/dev/dri rw
DeviceAllow=/dev/video* rw
DeviceAllow=char-input rw

# Variables d'environnement pour accélération matérielle
Environment=VLC_PLUGIN_PATH=/usr/lib/x86_64-linux-gnu/vlc/plugins
Environment=LIBVA_DRIVER_NAME=auto
Environment=VDPAU_DRIVER=auto

# Optimisations système
LimitNOFILE=65536
LimitMEMLOCK=infinity
LimitRTPRIO=95
LimitRTTIME=infinity
Nice=-10
IOSchedulingClass=1
IOSchedulingPriority=7

# Logging
StandardOutput=append:/opt/pisignage/logs/systemd-engine.log
StandardError=append:/opt/pisignage/logs/systemd-engine.log
SyslogIdentifier=pisignage-v4

# Sécurité
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/pisignage
PrivateTmp=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
MemoryDenyWriteExecute=no

[Install]
WantedBy=graphical-session.target
Also=nginx.service php8.2-fpm.service